# imageroot/systemd/user/mail-webhooks.service
#
# Copyright (C) 2025 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

[Unit]
Description=Mail Webhooks API server (mail-webhooks, %i)
After=network.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=%S/state/environment
EnvironmentFile=-%S/state/secrets/passwords.env
# Simple cleanup before start
ExecStartPre=-/usr/bin/podman stop --ignore %N
ExecStartPre=-/usr/bin/podman rm --ignore --force %N
ExecStart=/usr/bin/podman run \
    --detach \
    --conmon-pidfile=%t/%n.pid \
    --pid-file=%t/%n.pid \
    --cgroups=no-conmon \
    --replace --name=%N \
    --publish=127.0.0.1:${TCP_PORT}:8080 \
    --env-file=%S/state/environment \
    --env-file=-%S/state/secrets/passwords.env \
    ${MAIL_WEBHOOKS_IMAGE:-ghcr.io/nethserver/mail-webhooks-api:latest}
ExecStop=/usr/bin/podman stop --ignore --time=10 %N
ExecStopPost=/usr/bin/podman rm --ignore --force %N
PIDFile=%t/%n.pid
Type=forking
Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=default.target

[Install]
WantedBy=default.target
