# imageroot/systemd/user/mail-monitor.service
[Unit]
Description=Mail Webhooks Monitor Service
After=mail-webhooks.service
PartOf=mail-webhooks.service

[Service]
Type=simple
ExecStart=%i/bin/mail-monitor --mongodb-url=$MONGODB_URL
Restart=always
RestartSec=5
User=mail-webhooks%i
Group=mail-webhooks%i
Environment=PATH=%i/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%i/environment
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mail-webhooks-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/%i
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
