# imageroot/systemd/user/mail-webhooks-simple.service
#
# Copyright (C) 2025 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Simple systemd service for debugging

[Unit]
Description=Mail Webhooks API server - Simple Version
After=network.target

[Service]
EnvironmentFile=%S/state/environment
EnvironmentFile=-%S/state/secrets/passwords.env
ExecStartPre=-/usr/bin/podman stop mail-webhooks-simple
ExecStartPre=-/usr/bin/podman rm mail-webhooks-simple
ExecStart=/usr/bin/podman run \
    --name=mail-webhooks-simple \
    --publish=127.0.0.1:${TCP_PORT}:8080 \
    --env-file=%S/state/environment \
    --env-file=-%S/state/secrets/passwords.env \
    ${MAIL_WEBHOOKS_IMAGE:-ghcr.io/nethserver/mail-webhooks-api:latest}
ExecStop=/usr/bin/podman stop mail-webhooks-simple
ExecStopPost=/usr/bin/podman rm mail-webhooks-simple
Type=exec
Restart=on-failure
RestartSec=30
TimeoutStartSec=300

[Install]
WantedBy=default.target
